# supervisord.ini
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:hblink]
command=/usr/local/bin/python3 /opt/HBlink3/bridge.py
directory=/opt/HBlink3/
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true


[program:hbmonitor]
command=/usr/local/bin/python3 /opt/HBmonitor/monitor.py
direcory=/opt/HBmonitor
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:parrot]
command=/usr/local/bin/python3 /opt/HBlink3/playback.py -c /opt/HBlink3/playback.cfg
directory=/opt/HBlink3/
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true